# tests/unit/crontabs/crontab2
@hourly sh purge-cache.sh

* 0 * * * sh backup.sh